# hashbang/.config/systemd/user/dehydrated.timer
[Unit]
Description=Run dehydrated -c periodically

[Timer]
OnBootSec=15min
OnCalendar=weekly

[Install]
WantedBy=default.target
